# rootfs/etc/services.d/healthcheck/run
#!/bin/execlineb -P
/bin/export PATH /usr/bin:/bin

foreground { sleep 10 }
foreground { if -n -x { @curl@/bin/curl --connect-timeout 15 --max-time 19 -s -o /dev/null -f 127.0.0.1:${HTTPD_PORT}/phpinfo.php }  kill 1 }
