# rootfs/etc/services.d/healthcheck/run
#!/bin/execlineb -P
/bin/export PATH /usr/bin:/bin
importas -D 0 HTTPD_PORT HTTPD_PORT

foreground { sleep 10 }
foreground { echo executing phpinfo }

#foreground { if -n -t { @curl@/bin/curl --connect-timeout 15 --max-time 19 -s -o /dev/null -f 127.0.0.1:${HTTPD_PORT}/phpinfo.php }  s6-svc -r /run/services.d/httpd  }
foreground { if -n -t { @curl@/bin/curl --connect-timeout 15 --max-time 19 -s -f 127.0.0.1:${HTTPD_PORT}/phpinfo.php }  s6-svc -r /run/services.d/httpd  }
