# rootfs/etc/services.d/httpd/run
#!/bin/execlineb -P
/bin/export PATH /usr/bin:/bin
exec @apacheHttpd@/bin/httpd -D FOREGROUND -d @rootfs@/etc/httpd
