# rootfs/etc/services.d/syslog/run
#!/bin/execlineb -P
fdmove -c 2 1
exec -c
s6-ipcserver  -1 -- /dev/log
fdmove -c 1 2
ucspilogd

