# rootfs/etc/services.d/cleantmp/run
#!/bin/execlineb -P
/bin/export PATH /usr/bin:/bin

foreground { sleep 10 }
foreground { echo cleanup tmp }
foreground { @findutils@/bin/find /tmp/ -ctime +3 -delete }
foreground { sleep 7200 }

